# services/qhy-camera/pkg/90-rusty-photon-qhy.rules
# Rusty Photon QHYCCD camera access (installed by rusty-photon-qhy-camera).
#
# Enumerated QHY cameras (VID 1618) become group-accessible to plugdev
# instead of world-writable. This file deliberately sorts after the SDK's
# 85-qhyccd.rules (installed by rusty-photon-qhy-firmware-install) so it
# tightens that file's blanket MODE="0666" to 0660 + plugdev.
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", GROUP="plugdev", MODE="0660"
# High-resolution sensors need more than the usbfs default (16 MB) for
# bulk transfers; QHYCCD recommends 200 MB.
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", RUN+="/bin/sh -c 'echo 200 > /sys/module/usbcore/parameters/usbfs_memory_mb || true'"
