# services/phd2-guider/pkg/rusty-photon-phd2-guider.service
[Unit]
Description=Rusty Photon phd2-guider - guider service wrapping PHD2 (port 11130)
After=network.target

[Service]
Type=simple
# Bare invocation defaults to the `serve` subcommand. Runs with built-in
# defaults (PHD2 at localhost:4400, bind 127.0.0.1:11130) when no config
# exists at /etc/rusty-photon/phd2-guider.json. A missing or not-yet-started
# PHD2 is not fatal: the service serves 503 on /health and keeps retrying
# the connection, so no ConditionPathExists gate is needed.
ExecStart=/usr/bin/rusty-photon-phd2-guider
Restart=on-failure
RestartSec=5
User=rusty-photon
Group=rusty-photon
Environment=RUST_LOG=info
Environment=HOME=/var/lib/rusty-photon
WorkingDirectory=/var/lib/rusty-photon
StateDirectory=rusty-photon/phd2-guider

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/rusty-photon
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
UMask=0027

[Install]
WantedBy=multi-user.target
